# systemd/anticipatory.service
[Unit]
Description=Anticipatory Music Model Service
After=network.target

[Service]
Type=exec
WorkingDirectory=%h/src/music-models/services/anticipatory
ExecStart=%h/.local/bin/uv run python server.py
Restart=on-failure
RestartSec=10

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=anticipatory

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
