# systemd/audioldm2.service
[Unit]
Description=AudioLDM2 Text-to-Audio Service
After=network.target

[Service]
Type=exec
WorkingDirectory=%h/src/music-models/services/audioldm2
ExecStart=%h/.local/bin/uv run python server.py
Restart=on-failure
RestartSec=10

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=audioldm2

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
