# systemd/deepseek.service
[Unit]
Description=DeepSeek Tool Calling Service
After=network.target

[Service]
Type=exec
WorkingDirectory=%h/src/music-models/services/deepseek
ExecStart=%h/.local/bin/uv run python server.py
Restart=on-failure
RestartSec=10

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=deepseek

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
