# systemd/musicgen.service
[Unit]
Description=MusicGen Text-to-Music Service
After=network.target

[Service]
Type=exec
WorkingDirectory=%h/src/music-models/services/musicgen
ExecStart=%h/.local/bin/uv run python server.py
Restart=on-failure
RestartSec=10

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=musicgen

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
